# deploy/openfang.service
[Unit]
Description=OpenFang Agent OS Daemon
Documentation=https://github.com/openfang-ai/openfang
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openfang
Group=openfang
Slice=openfang.slice
ExecStart=/usr/local/bin/openfang start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment
EnvironmentFile=-/etc/openfang/env
WorkingDirectory=/var/lib/openfang

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/openfang
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
# Cgroup write access required for per-agent pids.max sandbox. Delegate=
# below restricts what we can touch to our own cgroup subtree.
ProtectControlGroups=no
RestrictSUIDSGID=true
MemoryDenyWriteExecute=false
RestrictRealtime=true

# Cgroup v2 controller delegation — lets the openfang user create per-agent
# sub-cgroups under /sys/fs/cgroup/openfang.slice/openfang.service and set
# pids.max / memory.max / cpu.max on them.
Delegate=pids cpu memory

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
