# deploy/openfang.slice
[Unit]
Description=Slice for OpenFang Agent OS workloads
Documentation=https://github.com/openfang-ai/openfang

[Slice]
CPUAccounting=yes
MemoryAccounting=yes
TasksAccounting=yes
